# data/aside-overlay.service
[Unit]
Description=aside - Wayland text overlay
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/aside-overlay
Restart=always
RestartSec=3
StartLimitIntervalSec=0

[Install]
WantedBy=default.target
