# data/aside-overlay.service
[Unit]
Description=aside - Wayland text overlay
After=graphical-session.target
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=%h/.local/bin/aside-overlay
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
