# data/aside-daemon.service
[Unit]
Description=aside - LLM desktop assistant daemon
After=graphical-session.target

[Service]
Type=simple
Environment=PYTHONUNBUFFERED=1
ExecStart=%h/.local/lib/aside/venv/bin/python3 -m aside.daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
